# rules/reference/linux/T1053.003/ref_cron_modification.yml
title: "[REF] Cron Job Modification via crontab"
id: c5e0f3b4-9d6a-4c1e-0a4f-3f7e8d9c0a1f
status: stable
description: |
  REFERENCE RULE — for bad-pool sample acceptance only. Not an evaluation rule.
  Fires on any execution of crontab (adding, listing, removing cron jobs),
  tagged with key=cron_mod.
author: svp-sandbox
date: 2026-06-14
tags:
    - attack.persistence
    - attack.t1053.003
logsource:
    product: linux
    service: auditd
detection:
    selection:
        key: 'cron_mod'
    condition: selection
level: informational
