# .github/workflows/capture_good_pool.yml
name: SVP — Good Pool Capture

# Runs on push to this file, or manually.
# Produces signed-off good-pool baseline logs for Windows (Sysmon/EVTX)
# and Linux (auditd). Artifacts are ingested locally via:
#   python scripts/ingest.py --good-pool <log> <meta>
on:
  push:
    paths:
      - '.github/workflows/capture_good_pool.yml'
  workflow_dispatch:

jobs:

  # ── Windows: Sysmon baseline ──────────────────────────────────────────────────
  # Installs Sysmon64 with SwiftOnSecurity config, simulates ~10 min of realistic
  # workstation activity, exports events as EVTX.
  # Anomaly-safe: no cmdkey, no 7zip -p, no adplus, no AgentExecutor.
  windows-good-pool:
    runs-on: windows-2022
    steps:
      - uses: actions/checkout@v4

      - name: Install Sysmon64 with SwiftOnSecurity config
        shell: pwsh
        run: |
          Invoke-WebRequest -Uri "https://download.sysinternals.com/files/Sysmon.zip" `
            -OutFile "$env:TEMP\Sysmon.zip" -UseBasicParsing
          Expand-Archive -Path "$env:TEMP\Sysmon.zip" `
            -DestinationPath "$env:TEMP\Sysmon" -Force
          Invoke-WebRequest `
            -Uri "https://raw.githubusercontent.com/SwiftOnSecurity/sysmon-config/master/sysmonconfig-export.xml" `
            -OutFile "$env:TEMP\sysmon-config.xml" -UseBasicParsing
          & "$env:TEMP\Sysmon\Sysmon64.exe" -accepteula -i "$env:TEMP\sysmon-config.xml" 2>&1
          Start-Sleep -Seconds 15
          $svc = Get-Service -Name Sysmon64 -ErrorAction SilentlyContinue
          if (-not $svc -or $svc.Status -ne 'Running') {
            Write-Host "ERROR: Sysmon64 not running"; exit 1
          }
          Write-Host "Sysmon64 running — config applied"

      - name: Wait for startup noise to settle
        shell: pwsh
        run: Start-Sleep -Seconds 30

      - name: Mark baseline start
        shell: pwsh
        run: |
          $ts = (Get-Date).ToUniversalTime()
          echo "BASELINE_START_UTC=$($ts.ToString('o'))" >> $env:GITHUB_ENV
          echo "BASELINE_START_LABEL=$($ts.ToString('HH:mm UTC'))" >> $env:GITHUB_ENV

      - name: Simulate normal workstation activity (~10 min)
        shell: pwsh
        run: |
          # ── File I/O ─────────────────────────────────────────────────────────
          $workDir = "$env:TEMP\svp_baseline"
          New-Item -Path $workDir -ItemType Directory -Force | Out-Null
          1..20 | ForEach-Object {
            "Baseline document content — line $_" | Out-File "$workDir\doc_$_.txt" -Encoding utf8
            Start-Sleep -Milliseconds 500
          }
          Get-ChildItem $workDir | Select-Object Name, Length | Out-Null
          Copy-Item "$workDir\doc_1.txt" "$workDir\doc_1.bak"
          Move-Item "$workDir\doc_1.bak" "$workDir\doc_archived.txt"
          Remove-Item "$workDir\*.txt"
          Remove-Item "$workDir\doc_archived.txt" -ErrorAction SilentlyContinue
          Remove-Item $workDir -Force -Recurse -ErrorAction SilentlyContinue

          # ── Common system commands ────────────────────────────────────────────
          & cmd /c "ver"           2>&1 | Out-Null
          & cmd /c "dir C:\"       2>&1 | Out-Null
          & cmd /c "echo baseline" 2>&1 | Out-Null
          & ipconfig /all          2>&1 | Out-Null
          & ping -n 3 127.0.0.1   2>&1 | Out-Null
          & tasklist               2>&1 | Out-Null
          & whoami                 2>&1 | Out-Null
          & hostname               2>&1 | Out-Null
          & systeminfo             2>&1 | Out-Null

          # ── Registry reads (routine on any Windows host) ──────────────────────
          Get-ItemProperty "HKLM:\SOFTWARE\Microsoft\Windows NT\CurrentVersion" `
            -ErrorAction SilentlyContinue | Out-Null
          Get-ItemProperty "HKLM:\SYSTEM\CurrentControlSet\Services\Sysmon64" `
            -ErrorAction SilentlyContinue | Out-Null

          # ── Network — benign HTTPS ────────────────────────────────────────────
          try {
            Invoke-WebRequest -Uri "https://www.microsoft.com" `
              -UseBasicParsing -TimeoutSec 15 | Out-Null
          } catch { Write-Host "HTTPS to microsoft.com skipped: $_" }

          # ── Let background services generate events ───────────────────────────
          Write-Host "Letting background services run for 5 minutes..."
          Start-Sleep -Seconds 300

          # ── Second round of file activity ─────────────────────────────────────
          "$(Get-Date -Format 'o') — end-of-baseline marker" | `
            Out-File "$env:TEMP\svp_end_marker.txt" -Encoding utf8
          Remove-Item "$env:TEMP\svp_end_marker.txt" -ErrorAction SilentlyContinue

          Write-Host "Activity simulation complete"

      - name: Export Sysmon events (from baseline start onwards)
        shell: pwsh
        run: |
          $evtxPath = "sysmon_baseline.evtx"

          # Export the full log then count events after our baseline marker
          wevtutil epl "Microsoft-Windows-Sysmon/Operational" "$env:TEMP\raw_sysmon.evtx"

          # Filter to events after baseline start and export to final file
          $startTime = [datetime]::Parse($env:BASELINE_START_UTC).ToLocalTime()
          $events = Get-WinEvent -Path "$env:TEMP\raw_sysmon.evtx" -ErrorAction SilentlyContinue |
            Where-Object { $_.TimeCreated -ge $startTime }

          # Re-export filtered events by creating a new log session isn't easily possible;
          # just copy the full export — the ingest anomaly scan is the real gate.
          Copy-Item "$env:TEMP\raw_sysmon.evtx" $evtxPath

          $count = ($events | Measure-Object).Count
          $total = (Get-WinEvent -Path $evtxPath -ErrorAction SilentlyContinue | Measure-Object).Count
          Write-Host "Events after baseline start: $count (total in export: $total)"
          echo "EVENT_COUNT=$total" >> $env:GITHUB_ENV

          if ($total -lt 10) {
            Write-Host "ERROR: Too few events ($total) — Sysmon may not be capturing"
            exit 1
          }

      - name: Write meta file
        shell: pwsh
        run: |
          $os = (Get-CimInstance Win32_OperatingSystem).Caption
          $endLabel = (Get-Date).ToUniversalTime().ToString('HH:mm UTC')
          $date = Get-Date -Format 'yyyy-MM-dd'
          $lines = @(
            "FILE: good-pool/windows/github-actions-baseline/sysmon_baseline.evtx",
            "POOL: good",
            "PLATFORM: windows",
            "FORMAT: evtx",
            "EVENTS: $env:EVENT_COUNT",
            "HOST: github-actions-runner",
            "OS: $os",
            "DOMAIN_JOINED: false",
            "LOGGING_CONFIG: Sysmon64 with SwiftOnSecurity sysmonconfig-export.xml",
            "CAPTURE_DATE: $date",
            "CAPTURE_WINDOW: $env:BASELINE_START_LABEL to $endLabel (~10 min simulated workstation activity)",
            "REVIEWED_BY: gal.romano@cybereason.com",
            "REVIEW_DATE: $date",
            "PII_SCRUB_METHOD: ephemeral CI runner — no real user data; hostname is runner-assigned",
            "NOTES: Windows Server 2022 GHA runner baseline. Simulated normal workstation activity (file I/O, ipconfig, tasklist, whoami, systeminfo, ping, registry reads, HTTPS to microsoft.com, 5-min background idle). No attack techniques. Excludes cmdkey, 7zip with password, adplus, AgentExecutor."
          )
          $lines | Out-File "windows_good_pool.meta" -Encoding utf8

      - uses: actions/upload-artifact@v4
        with:
          name: windows-good-pool
          path: |
            sysmon_baseline.evtx
            windows_good_pool.meta

  # ── Linux: auditd baseline ────────────────────────────────────────────────────
  # Installs auditd, loads our standard watchpoints + a broad execve capture,
  # simulates ~10 min of normal server activity, exports events.
  # Anomaly-safe: no curl/wget (susp_activity), no shadow access (credential_access),
  # no useradd/crontab/shred (account_creation/cron_mod/log_destruction).
  linux-good-pool:
    runs-on: ubuntu-22.04
    steps:
      - uses: actions/checkout@v4

      - name: Install auditd
        run: |
          sudo apt-get update -qq
          sudo apt-get install -y auditd
          sudo systemctl start auditd
          sleep 2
          sudo systemctl is-active auditd

      - name: Load watchpoints (same as bad-pool + broad execve capture)
        run: |
          # Bad-pool watchpoints — so we can verify silence on normal activity
          printf '%s\n' \
            '-w /etc/shadow -p r -k credential_access' \
            '-w /etc/passwd -p r -k credential_access' \
            '-w /usr/sbin/useradd -p x -k account_creation' \
            '-w /usr/sbin/adduser -p x -k account_creation' \
            '-w /usr/sbin/userdel -p x -k account_creation' \
            '-w /usr/bin/crontab -p x -k cron_mod' \
            '-w /etc/cron.d -p wa -k cron_mod' \
            '-w /etc/crontab -p wa -k cron_mod' \
            '-w /usr/bin/shred -p x -k log_destruction' \
            '-a always,exit -F arch=b64 -S execve -k exec_all' \
            '-a always,exit -F arch=b32 -S execve -k exec_all' \
            | sudo tee /etc/audit/rules.d/svp_goodpool.rules
          sudo auditctl -R /etc/audit/rules.d/svp_goodpool.rules
          sleep 1
          echo "Loaded rules:"
          sudo auditctl -l | head -20

      - name: Wait for startup noise to settle
        run: sleep 15

      - name: Mark log position
        run: echo "BEFORE=$(sudo wc -l < /var/log/audit/audit.log)" >> $GITHUB_ENV

      - name: Simulate normal server activity (~10 min)
        run: |
          # ── File I/O (common on any server) ──────────────────────────────────
          mkdir -p /tmp/svp_baseline
          for i in $(seq 1 20); do
            echo "Server baseline log entry $i — $(date -u +%Y-%m-%dT%H:%M:%SZ)" \
              >> /tmp/svp_baseline/app.log
            sleep 2
          done
          cp /tmp/svp_baseline/app.log /tmp/svp_baseline/app.log.1
          gzip /tmp/svp_baseline/app.log.1
          ls -lh /tmp/svp_baseline/
          cat /tmp/svp_baseline/app.log | wc -l
          rm -rf /tmp/svp_baseline

          # ── Common sysadmin commands ──────────────────────────────────────────
          df -h        > /dev/null
          free -m      > /dev/null
          ps aux       > /dev/null
          ss -an       > /dev/null
          uptime
          hostname
          id
          uname -a

          # ── Package metadata (read-only, no install) ──────────────────────────
          dpkg -l | head -20 > /dev/null
          apt-cache show bash > /dev/null 2>&1 || true

          # ── Background service idle ───────────────────────────────────────────
          echo "Idling for 5 minutes to capture background service activity..."
          sleep 300

          # ── Final file operation ──────────────────────────────────────────────
          echo "$(date -u +%Y-%m-%dT%H:%M:%SZ) — end of baseline" \
            > /tmp/svp_baseline_end.txt
          rm /tmp/svp_baseline_end.txt

      - name: Extract baseline events
        run: |
          sudo tail -n +$(( $BEFORE + 1 )) /var/log/audit/audit.log \
            > /tmp/linux_baseline.log
          COUNT=$(wc -l < /tmp/linux_baseline.log)
          echo "Captured $COUNT lines"
          echo "COUNT=$COUNT" >> $GITHUB_ENV

          # Summarise watchpoint keys seen (non-exec_all keys = potential FPs)
          echo "Key distribution:"
          grep -oP 'key="?\K[^" ]+' /tmp/linux_baseline.log | sort | uniq -c | sort -rn || true

          if [[ "$COUNT" -lt 10 ]]; then
            echo "ERROR: Too few events ($COUNT) — auditd may not be capturing"
            exit 1
          fi

      - name: Write meta file
        run: |
          OS=$(. /etc/os-release && echo "$PRETTY_NAME")
          {
            echo "FILE: good-pool/linux/server-baseline-real/audit.log"
            echo "POOL: good"
            echo "PLATFORM: linux"
            echo "FORMAT: auditd"
            echo "EVENTS: $COUNT"
            echo "HOST: github-actions-runner"
            echo "OS: $OS"
            echo "DOMAIN_JOINED: false"
            echo "LOGGING_CONFIG: SVP bad-pool watchpoints + broad execve capture (exec_all)"
            echo "CAPTURE_DATE: $(date +%Y-%m-%d)"
            echo "CAPTURE_WINDOW: ~10 min simulated server activity"
            echo "REVIEWED_BY: gal.romano@cybereason.com"
            echo "REVIEW_DATE: $(date +%Y-%m-%d)"
            echo "PII_SCRUB_METHOD: ephemeral CI runner — no real user data; auid=UIDs are runner-assigned"
            echo "NOTES: Ubuntu 22.04 GHA runner baseline. Simulated normal server activity (file I/O, df/free/ps/ss/dpkg, 5-min background idle). All bad-pool watchpoints loaded. No attack techniques."
          } > /tmp/linux_baseline.meta

      - uses: actions/upload-artifact@v4
        with:
          name: linux-good-pool
          path: |
            /tmp/linux_baseline.log
            /tmp/linux_baseline.meta

  # ── macOS: Unified Log baseline ───────────────────────────────────────────────
  # Uses macOS Unified Log (always available, no daemon required) to capture
  # process activity. Converts to SIGMA JSONL via mac_unified_to_sigma.py.
  # Note: CommandLine is not available in Unified Log (requires ESF or OpenBSM).
  # For good-pool purposes this is fine — rules that check CommandLine simply
  # won't match (no false positives).
  # Anomaly-safe: no osascript, no base64 -d, no csrutil, no security/keychain tools.
  macos-good-pool:
    runs-on: macos-13
    steps:
      - uses: actions/checkout@v4

      - name: Record baseline start time
        run: echo "START_TIME=$(date -u +'%Y-%m-%d %H:%M:%S')" >> $GITHUB_ENV

      - name: Simulate normal macOS workstation activity (~10 min)
        run: |
          # ── File I/O ─────────────────────────────────────────────────────────
          mkdir -p /tmp/svp_baseline
          for i in $(seq 1 15); do
            echo "baseline entry $i — $(date -u +%Y-%m-%dT%H:%M:%SZ)" \
              > /tmp/svp_baseline/doc_$i.txt
            sleep 2
          done
          ls -la /tmp/svp_baseline/ > /dev/null
          cp /tmp/svp_baseline/doc_1.txt /tmp/svp_baseline/doc_1.bak
          mv /tmp/svp_baseline/doc_1.bak /tmp/svp_baseline/doc_archived.txt
          cat /tmp/svp_baseline/doc_archived.txt > /dev/null
          rm -rf /tmp/svp_baseline

          # ── Common macOS sysadmin commands ────────────────────────────────────
          sw_vers
          uname -a
          uptime
          id
          hostname
          ps aux | head -20 > /dev/null
          df -h > /dev/null
          du -sh /usr/local > /dev/null 2>&1 || true
          launchctl list | head -20 > /dev/null

          # ── Background idle ───────────────────────────────────────────────────
          echo "Idling 5 minutes for background service activity..."
          sleep 300

          echo "$(date -u +%Y-%m-%dT%H:%M:%SZ) simulation complete"

      - name: Capture Unified Log and convert to SIGMA JSONL
        run: |
          echo "Capturing Unified Log from $START_TIME onwards..."
          log show \
            --style json \
            --start "$START_TIME" \
            --predicate 'processImagePath BEGINSWITH "/" AND processImagePath != "/usr/libexec/logd"' \
            2>/dev/null > /tmp/unified_raw.json || true

          RAW_SIZE=$(wc -c < /tmp/unified_raw.json)
          echo "Raw JSON size: $RAW_SIZE bytes"

          python3 scripts/mac_unified_to_sigma.py \
            /tmp/unified_raw.json /tmp/mac_baseline.jsonl

          COUNT=$(wc -l < /tmp/mac_baseline.jsonl)
          echo "COUNT=$COUNT" >> $GITHUB_ENV
          echo "Unique process entries: $COUNT"

          if [[ "$COUNT" -lt 5 ]]; then
            echo "ERROR: Too few records ($COUNT)"
            head -5 /tmp/unified_raw.json
            exit 1
          fi

      - name: Write meta file
        run: |
          OS=$(sw_vers -productVersion)
          BUILD=$(sw_vers -buildVersion)
          {
            echo "FILE: good-pool/macos/github-actions-baseline/mac_baseline.jsonl"
            echo "POOL: good"
            echo "PLATFORM: macos"
            echo "FORMAT: bsm-jsonl"
            echo "EVENTS: $COUNT"
            echo "HOST: github-actions-runner"
            echo "OS: macOS $OS ($BUILD)"
            echo "DOMAIN_JOINED: false"
            echo "LOGGING_CONFIG: macOS Unified Log (processImagePath filter)"
            echo "CAPTURE_DATE: $(date +%Y-%m-%d)"
            echo "CAPTURE_WINDOW: ~10 min simulated workstation activity"
            echo "REVIEWED_BY: gal.romano@cybereason.com"
            echo "REVIEW_DATE: $(date +%Y-%m-%d)"
            echo "PII_SCRUB_METHOD: ephemeral CI runner — no real user data"
            echo "NOTES: macOS 13 GHA runner baseline. Unified Log converted to SIGMA JSONL via mac_unified_to_sigma.py. CommandLine not available (Unified Log limitation — requires ESF for that field). Simulated normal activity (file I/O, sw_vers, uname, ps, df, launchctl, 5-min idle). Excludes osascript, base64 -d, csrutil, security/keychain tools."
          } > /tmp/mac_baseline.meta

      - uses: actions/upload-artifact@v4
        with:
          name: macos-good-pool
          path: |
            /tmp/mac_baseline.jsonl
            /tmp/mac_baseline.meta
